# crontab
0 3 1 * * root sh /scripts/run_certbot.sh
